# docker/crontab
*/5 * * * * root /usr/local/bin/search_maint.sh > /dev/null 2>&1
